# systemd/user/kakaotalk-notify.service
[Unit]
Description=KakaoTalk notify (inotify -> notify-send)
After=graphical-session.target dunst.service
Wants=graphical-session.target dunst.service

[Service]
Type=simple
Environment=WINEPREFIX=%h/.local/share/kakaotalk
ExecStartPre=/bin/sleep 5
ExecStart=%h/.scripts/kakaotalk-notify.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target

